# redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/sbin/redis-server /etc/redis.conf

[program:sshd]
command=/usr/sbin/sshd -D
